# roles/zfs/templates/zpool-scrub@.timer.j2
[Timer]
Unit=%N.service
OnCalendar={{ item.value.scrub_schedule }}

[Install]
WantedBy=timers.target
