# roles/common/templates/item.mount
[Mount]
What={{ item.value.device }}
Where={{ item.value.mount_point }}
