# roles/common/templates/item.automount
[Automount]
Where={{ item.value.mount_point }}
TimeoutIdleSec={{ item.value.timeout_idle }}

[Install]
WantedBy=multi-user.target
