# roles/restic/templates/restic.timer.j2
[Timer]
Unit=%N.service
OnCalendar={{ restic_schedule }}

[Install]
WantedBy=timers.target
