# roles/common/files/ssh-pull-keys.service
[Service]
Type=oneshot
EnvironmentFile=/etc/default/%N
ExecStart=ssh-pull-keys
