# roles/common/files/ssh-pull-keys.timer
[Timer]
Unit=%N.service
OnCalendar=*:0

[Install]
WantedBy=timers.target
