# roles/services/templates/link-services-crontabs.path.j2
[Path]
PathModified={{ services_root }}

[Install]
WantedBy=multi-user.target
