# roles/apt/templates/apt-upgrade.timer.j2
[Timer]
Unit=%N.service
OnCalendar={{ apt_upgrade_schedule }}

[Install]
WantedBy=timers.target
