# roles/heartbeat/templates/heartbeat.service.j2
[Service]
EnvironmentFile=/etc/default/%N
ExecStart=sh -c 'while echo "$(curl -s https://hc-ping.com/${HEALTHCHECKS_UUID})"; do sleep 60; done'

[Install]
WantedBy=multi-user.target
